# system_files/etc/udev/rules.d/81-wol.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="enp37s0", RUN+="/usr/bin/ethtool -s enp37s0 wol g"
